# article-saver.service
[Unit]
Description=Article Saver - Lightweight Article Archiving System
After=network.target

[Service]
Type=simple
User=gantl
WorkingDirectory=/home/gantl/article-saver
Environment="PATH=/home/gantl/article-saver/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=/home/gantl/article-saver/.env
ExecStart=/home/gantl/article-saver/venv/bin/python3 backend/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
